# etc/services.d/logmonitor/run
#!/usr/bin/with-contenv sh

while true
do
    echo "DISCORD_WEBHOOK=${DISCORD_WEBHOOK}" > /.envvars
    sleep 3600
done

